# supervisord.conf
[supervisord]
nodaemon=true

[include]
files = /etc/supervisor/conf.d/*.conf

[program:uwsgi]
autorestart=true
command=/var/www/myproject/uwsgi/run.sh
priority=300
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
